# systemd/gunicorn.socket
[Unit]
Description=Gunicorn HTTP server socket

[Socket]
ListenStream=/run/gunicorn/socket

[Install]
WantedBy=sockets.target